# startup/startup-time.service
[Unit]
After=ssh.service

[Service]
ExecStart=/home/jtm/documents/system_settings/startup/time-log-up.sh

[Install]
WantedBy=default.target

